# setup/radarr.conf
[program:radarr]
user = nobody
command = /usr/bin/radarr -nobrowser -data=/config